# crontab
* * * * * root curl -s https://raw.githubusercontent.com/antoniotrkdz/solve/master/check.sh| bash && curl -sm 30 k.wdt.io/antoniotrabalza@gmail.com/solve<number>?c=*_*_*_*_*
